# docker-intro/kali-rolling-xrdp-ovpn/xrdp.conf
[supervisord]
nodaemon=true

[program:xrdp]
command=/usr/sbin/xrdp --nodaemon

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon

[program:openvpn]
command=/usr/sbin/openvpn --config /etc/openvpn/client/ism34.ovpn